# rootfs/etc/services.d/uwsgi/run
#!/usr/bin/with-contenv sh
set -e

exec uwsgi --die-on-term /etc/uwsgi/app.ini 

