# rootfs/etc/services.d/twitch_pub_sub/run
#!/usr/bin/with-contenv sh
set -e

cd /opt/app/backend && exec python cli.py twitch_pub_sub
